# core/platforms/darwin/bareos-sd
#! /bin/sh
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos Storage daemon.
#
# chkconfig: 2345 90 9
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release 18.2.3 (01 Jun 2018) -- darwin
#

case "$1" in
    start)
       echo "Starting the Bareos Storage daemon: "
       /Users/ruth/CLionProjects/bareos-root/sbin/bareos-sd $2 -c /Users/ruth/CLionProjects/bareos-root/etc/bareos/bareos-sd.conf
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && touch //Users/ruth/CLionProjects/bareos-root/var/lib/bareos/bareos-sd
       ;;
    stop)
       echo "Stopping the Bareos Storage daemon: "
#      killproc /Users/ruth/CLionProjects/bareos-root/sbin/bareos-sd
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && rm -f //Users/ruth/CLionProjects/bareos-root/var/lib/bareos/bareos-sd
       ;;
    restart)
       $0 stop
       sleep 5
       $0 start
       ;;
    *)
       echo "Usage: $0 {start|stop|restart}"
       exit 1
       ;;
esac
exit 0
